# hw8_prometheus/prometheus.ini
[program:prometheus]


command=/home/user27/hw8_prometheus/prometheus-2.3.2.linux-amd64/prometheus
        --config.file=/home/user27/hw8_prometheus/prometheus-2.3.2.linux-amd64/prometheus.yml
        --web.console.libraries=/home/user27/hw8_prometheus/prometheus-2.3.2.linux-amd64/console_libraries
        --web.console.templates=/home/user27/hw8_prometheus/prometheus-2.3.2.linux-amd64/consoles

stopsignal=QUIT
autostart=true
startretries=10
startsecs=0
stopwaitsecs=10
stopasgroup=true

stdout_logfile=/var/log/prometheus/prometheus_stdout.log
stdout_logfile_maxbytes=60MB
stdout_logfile_backups=4
stdout_capture_maxbytes=4MB

stderr_logfile=/var/log/prometheus/prometheus_stderr.log
stderr_logfile_maxbytes=60MB
stderr_logfile_backups=4
stderr_capture_maxbytes=4MB
