# hw8_prometheus/uwsgi_exporter.ini
[program:uwsgi_exporter]


command=/home/user27/hw8_prometheus/uwsgi_exporter-0.7.0.linux-amd64/uwsgi_exporter
        --stats.uri="http://127.0.0.1:1717"


stopsignal=QUIT
autostart=true
startretries=10
startsecs=0
stopwaitsecs=10
stopasgroup=true

stdout_logfile=/var/log/uwsgi_exporter/uwsgi_exporter_stdout.log
stdout_logfile_maxbytes=60MB
stdout_logfile_backups=4
stdout_capture_maxbytes=4MB

stderr_logfile=/var/log/uwsgi_exporter/uwsgi_exporter_stderr.log
stderr_logfile_maxbytes=60MB
stderr_logfile_backups=4
stderr_capture_maxbytes=4MB
